# scripts/obsidian-sync.timer
[Unit]
Description=Run Obsidian vault sync every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
